# api_qbwc_zoho_backend/project_api/crontab
* * * * * /usr/local/bin/python /app/manage.py runcrons >> /var/log/cron.log 2>&1
